# bot/conf/bot_conf.conf
[program:bot]
command=/bot/venv/bin/python /bot/bots/bot_father.py
environment=PATH=/bot/venv/bin
stdout_logfile=/bot/error.log
